# perl-modules/Mail-SpamAssassin/spamd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='mail/spamassassin/spamd' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='network' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <method_context/>
    <!-- log to stderr, ie. smf logfile. even logging to syslog requires spamd
         to get HUPed on syslog restart according to spamd(1). additionally
         using -d to daemonize fails because the child is unable to create a
         pidfile and the parent bails; just use & instead -->
    <exec_method name='start' type='method' exec='/opt/niksula/bin/spamd -s stderr &amp;' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
    <property_group name='application' type='application'/>
    <property_group name='startd' type='framework'>
      <propval name='ignore_error' type='astring' value='core,signal'/>
      <propval name='duration' type='astring' value='contract'/>
    </property_group>
    <stability value='Evolving'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>spamassassin spamd</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
